# etc/systemd/system/systemd-logind.service.d/override.conf
# Required to make systemd-logind work
[Service]
ProtectHostname=no
