# etc/systemd/system/console-getty.service.d/override.conf
# makes sure that the gnomeshell user gets logged in automatically
[Service]
ExecStart=
ExecStart=-/usr/sbin/agetty --autologin gnomeshell --noclear --keep-baud console 115200,38400,9600 $TERM
